# chef/cookbooks/baseconfig/files/default/rc.local
#!/bin/sh
sleep 10
sudo /usr/local/bin/uwsgi --ini /home/vagrant/projectCityEvents/wsgi/uwsgi.ini --daemonize /var/log/projectCityEvents.log
